# 8/sparky-desktop-data/manokwari/skel/.config/autostart/sparky-polkit.desktop
[Desktop Entry]
Comment=Sparky Polkit
Encoding=UTF-8
Exec=/usr/bin/sparky-polkit
Name=Sparky Polkit
StartupNotify=true
Terminal=false
